# scripts/phab-http-reports.service
[Unit]
Description=phab-http reports

[Service]
Type=oneshot
ExecStart=/opt/epiphyte/phab-http/reports.sh
