# scripts/phab-http-reports.service
[Unit]
Description=phab-http reports

[Service]
Type=oneshot
ExecStart=/usr/bin/phab-http-reports
